# roles/laptop/files/systemd/disable-xhc-wakeup.service
[Unit]
Description=Disable XHC wakeup

[Service]
ExecStart=/bin/bash -c "echo XHC >> /proc/acpi/wakeup"

[Install]
WantedBy=multi-user.target
